# .claude/launchd/com.begifted.sentry-autofix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.begifted.sentry-autofix</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>"$HOME/development/begifted-native/.claude/scripts/sentry-autofix-cron"</string>
	</array>
	<key>StartCalendarInterval</key>
	<array>
		<dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>7</integer></dict>
		<dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>7</integer></dict>
		<dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>7</integer></dict>
		<dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>7</integer></dict>
		<dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>7</integer></dict>
	</array>
	<key>StandardOutPath</key>
	<string>/Users/aaronmullan/Library/Logs/begifted-autofix/launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/aaronmullan/Library/Logs/begifted-autofix/launchd.log</string>
</dict>
</plist>
